# deploy/bestme.service
[Unit]
Description=BestMe API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bestme
Environment=ADDR=127.0.0.1:8090
Environment=DB_PATH=/opt/bestme/data/bestme.db
Environment=WEB_DIST=/opt/bestme/web
EnvironmentFile=-/opt/bestme/bestme.env
ExecStart=/opt/bestme/bin/bestme
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
